# node-controller/phantom-node-controller.service
[Unit]
Description=Phantom Node Controller for FPTN
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/phantom-node-controller.env
ExecStart=/bin/sh -lc 'if [ -x /opt/phantom-node-controller/.venv/bin/python ]; then exec /opt/phantom-node-controller/.venv/bin/python /opt/phantom-node-controller/agent.py; else exec /usr/bin/python3 /opt/phantom-node-controller/agent.py; fi'
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
