# deploy/phantom-control-plane.service
[Unit]
Description=Phantom Control Plane
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=phantom
Group=phantom
WorkingDirectory=/opt/phantom-control-plane
EnvironmentFile=/etc/phantom-control-plane.env
ExecStart=/bin/sh -lc 'exec /opt/phantom-control-plane/.venv/bin/python -m uvicorn app.main:app --host "${PANEL_HOST:-0.0.0.0}" --port "${PANEL_PORT:-8000}"'
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
